# crontab
30 */3 * * *  ${HOME}/.cron-script/run.sh homebrew-update
45 */3 * * *  ${HOME}/.cron-script/run.sh asdf-update
* * * * *  ${HOME}/.cron-script/run.sh backup-btt
